# deploy/supervisord.conf
; ONC center 单容器双进程：gunicorn（Flask 应用/:8080）+ caddy（web 前置 80/443）
; 两者 autorestart；日志直写容器 stdout/stderr，docker logs 可见。
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:gunicorn]
command=gunicorn -b 0.0.0.0:8080 -w 2 --timeout 60 app:app
directory=/app
autostart=true
autorestart=true
startsecs=3
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
directory=/app
autostart=true
autorestart=true
startsecs=3
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
